# admin_only/bco_api.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/var/run/bco_api.sock

[Install]
WantedBy=sockets.target
